# siduction-btrfs-0.1.2/siduction_btrfs.path
# /usr/lib/systemd/system/siduction_btrfs.path
[Unit]
Description=Part of siduction-btrfs
Description=Monitors changes of /boot/grub/grub-btrfs.cfg
Documentation=https://github.com/siduction/siduction-btrfs
After=siduction_btrfs.service

[Path]
PathModified=/boot/grub/grub-btrfs.cfg

[Install]
WantedBy=multi-user.target
