# siduction-btrfs-0.1.0/debian/postinst
#! /bin/bash
#
# postinst script for #PACKAGE#
#

set -e

systemctl enable siduction_btrfs.timer siduction_btrfs.path

#DEBHELPER#

exit 0
