# siduction-btrfs-0.1.0/debian/siduction-btrfs/DEBIAN/postinst
#! /bin/bash
#
# Name: postinst
# Part of siduction-btrfs

systemctl enable --now siduction_btrfs.timer
systemctl enable --now siduction_btrfs.path

exit 0
